# systemd/hermes-monitor.service.d/on-failure-alert.conf
[Unit]
OnFailure=hermes-monitor-alert.service
